# services/xdm/run
#!/bin/bash
exec /sbin/setuser abc "/usr/bin/python /var/XDM/XDM.py -b /config/data -n" 